# scripts/crontab
0 */12 * * * /scripts/renew-certs.sh >> /var/log/cron.log 2>&1